# startup/upstart/jane.conf
description "jane"
author "Project Jane"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

console log

exec /usr/bin/jane
