# startup/systemd/jane.service
[Unit]
Description=Hex Bot

[Service]
Type=simple
ExecStart=/usr/local/bin/hex
Restart=on-failure
