# startup/systemd/jane.service
[Unit]
Description=Jane Bot

[Service]
Type=simple
ExecStart=/usr/local/bin/jane
Restart=on-failure
